# deploy/ari-mm-update.service
# /etc/systemd/system/ari-mm-update.service
[Unit]
Description=Pull latest ARI branch and restart the manager
[Service]
Type=oneshot
User=ariapp
ExecStart=/opt/ari/repo/metadata-manager_v2/deploy/update.sh
