# deploy/ari-mm-update.service
# /etc/systemd/system/ari-mm-update.service
[Unit]
Description=Pull latest ARI Metadata Manager app code and restart the manager
[Service]
Type=oneshot
User=ariapp
ExecStart=/opt/ari/ari-metadata-manager/deploy/update.sh
