# deploy/ari-mm-ontology-update.service
# /etc/systemd/system/ari-mm-ontology-update.service
[Unit]
Description=Fetch latest ARI ontology from GitHub and restart the manager if changed
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ariapp
ExecStart=/opt/ari/ari-metadata-manager/deploy/update-ontology.sh