# deploy/ari-mm-update.timer
# /etc/systemd/system/ari-mm-update.timer  ->  enable with: systemctl enable --now ari-mm-update.timer
[Unit]
Description=Periodically refresh ARI Metadata Manager app code
[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
[Install]
WantedBy=timers.target
